# cron/crontab
# =============================================================================
# SERVISTECH ERP V4.0 - Cron Jobs Schedule
# America/Caracas timezone (UTC-4)
# =============================================================================

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# Log all output
* * * * * echo "Cron is running at $(date)" >> /var/log/cron/cron.log 2>&1

# =============================================================================
# BCV Rate Scraping - 8:00 AM and 1:00 PM Venezuela time
# =============================================================================
0 8 * * 1-6 /app/scripts/scrape-bcv.sh >> /var/log/cron/bcv.log 2>&1
0 13 * * 1-6 /app/scripts/scrape-bcv.sh >> /var/log/cron/bcv.log 2>&1

# =============================================================================
# Binance USDT Rate Update - Every 30 minutes during business hours
# =============================================================================
*/30 7-19 * * 1-6 /app/scripts/update-binance.sh >> /var/log/cron/binance.log 2>&1

# =============================================================================
# Database Backups - Daily at 2:00 AM
# =============================================================================
0 2 * * * /app/scripts/backup-database.sh >> /var/log/cron/backup.log 2>&1

# =============================================================================
# Weekly Full Backup - Sundays at 3:00 AM
# =============================================================================
0 3 * * 0 /app/scripts/backup-full.sh >> /var/log/cron/backup-full.log 2>&1

# =============================================================================
# Clean old backups - Keep last 30 days
# =============================================================================
0 4 * * * /app/scripts/cleanup-backups.sh >> /var/log/cron/cleanup.log 2>&1

# =============================================================================
# Commission Calculation - Daily at 11:59 PM
# =============================================================================
59 23 * * * /app/scripts/calculate-commissions.sh >> /var/log/cron/commissions.log 2>&1

# =============================================================================
# Health Check Report - Daily at 9:00 AM
# =============================================================================
0 9 * * * /app/scripts/health-report.sh >> /var/log/cron/health.log 2>&1
